# ARCHIVE/ANSIBLE-COLLECTION/roles/rss2email/templates/etc_systemd_system_rss2email.timer.j2
[Unit]
Description=Run rss2email
After=network-online.target

[Install]
WantedBy=timers.target

[Timer]
Unit=rss2email.service
Persistent=true
OnCalendar=*-*-* *:00:11
OnBootSec=1s
